# scripts/mykey-proxy-daemon.service
[Unit]
Description=MyKey Proxy Daemon
Documentation=https://github.com/JamesFromFL/mykey-proxy
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/mykey-proxy-daemon
Environment=DBUS_SYSTEM_BUS_ADDRESS=unix:path=/run/dbus/system_bus_socket
Restart=on-failure
RestartSec=5s

# Run as a dedicated unprivileged system user
User=mykey-proxy

# Systemd hardening
NoNewPrivileges=yes
ProtectSystem=strict
MemoryDenyWriteExecute=yes
ReadWritePaths=/etc/mykey-proxy
ReadOnlyPaths=/etc
ReadWritePaths=/tmp
PrivateTmp=yes
RestrictRealtime=yes
SystemCallFilter=@system-service

# Allow the daemon to lock memory (needed for mlock on session tokens)
LimitMEMLOCK=65536

[Install]
WantedBy=default.target
